# systemd_user/low_battery_notif.service
[Unit]
Description=check-battery

[Service]
ExecStart=/home/phoenix/.scripts/low_battery_notif
