# systemd_user/low_battery_notif.service
[Unit]
Description=check-battery

[Service]
ExecStart=##absolute path here##/.scripts/low_battery_notif
